# ina228-monitor.service
[Unit]
Description=INA228 Multi-Channel Current Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yog
WorkingDirectory=/home/yog/Projects/power-monitor-apollo-fleet
ExecStart=/home/yog/Projects/power-monitor-apollo-fleet/.venv/bin/python3 /home/yog/Projects/power-monitor-apollo-fleet/monitor.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ina228-monitor

[Install]
WantedBy=multi-user.target
